# practise/celery_example1/work2_supervisor.ini
[program:work2]
command=venv/bin/celery -A practise.celery_example1.worker2 worker --loglevel=info 
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
startretries=2                ;
stopsignal=TERM               ; signal used to kill process (default TERM)
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
directory=/root/eagleeye/
autostart=true
stdout_logfile =/var/log/work2.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile =/var/log/work2.err
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10

